# crontab
45 8 * * * sudo docker stop bot && sudo docker start bot
