# crontab
45 7 * * * sudo docker stop bot && sudo docker start bot
